# root/config/cron/periodic/daily/tls-config-update
#!/bin/bash 

for generator in /app/letsencrypt/generated/letsencrypt/generate-*; do
    chmod +x "$generator"
    "$generator"
done